# chezmoi/private_dot_config/autostart/executable_kando.desktop
[Desktop Entry]
Type=Application
Name=Kando
Comment=Iniciar Kando no login
Exec=/usr/bin/kando
Terminal=false
X-GNOME-Autostart-enabled=true
