# chezmoi/private_dot_config/autostart/executable_albert.desktop
[Desktop Entry]
Type=Application
Name=Albert
Comment=Iniciar Albert no login
Exec=/usr/bin/albert
Terminal=false
X-GNOME-Autostart-enabled=true
