# chezmoi/private_dot_config/autostart/xkb-cedilha.desktop
[Desktop Entry]
Type=Application
Name=XKB Cedilha
Exec=/bin/sh -c "$HOME/.xkb-cedilha.sh"
X-GNOME-Autostart-enabled=true
